# crontab
SHELL=/bin/bash
PATH=/opt/news_app/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Backup PostgreSQL daily (02:00 UTC)
0 2 * * * flock -n /tmp/news_app_backup.lock /bin/bash -lc '/opt/news_app/scripts/backup_database.sh' >> /var/log/news_app/backup.log 2>&1

# Run scrapers every 15 minutes
*/15 * * * * flock -n /tmp/news_app_scrapers.lock /bin/bash -lc 'cd /opt/news_app && /opt/news_app/.venv/bin/python scripts/run_scrapers.py --show-stats --scrapers HackerNews Reddit Substack Techmeme Podcast Atom' >> /var/log/news_app/scrapers-cron.log 2>&1

# Run all scheduled Twitter/X work every 15 minutes
*/15 * * * * flock -n /tmp/news_app_twitter.lock /bin/bash -lc 'cd /opt/news_app && /opt/news_app/.venv/bin/python scripts/run_twitter.py' >> /var/log/news_app/twitter.log 2>&1

# Run news-native digest enqueue every 15 minutes

# Run feed discovery weekly (Mondays at 03:00 UTC)
0 3 * * 1 flock -n /tmp/news_app_feed_discovery.lock /bin/bash -lc 'cd /opt/news_app && /opt/news_app/.venv/bin/python scripts/run_feed_discovery.py' >> /var/log/news_app/feed-discovery-cron.log 2>&1
