# supervisor.conf
# News App Supervisor Configuration
# Copy to: /etc/supervisor/conf.d/news_app.conf
# Then run: supervisorctl reread && supervisorctl update

[program:news_app_server]
command=/bin/bash -lc "/opt/news_app/scripts/start_server.sh"
directory=/opt/news_app
user=newsapp
autostart=true
autorestart=true
stopsignal=TERM
startretries=3
stdout_logfile=/var/log/news_app/server.log
stderr_logfile=/var/log/news_app/server.err.log
redirect_stderr=false
environment=ENVIRONMENT="production"

[program:news_app_workers_content]
command=/bin/bash -lc "/opt/news_app/.venv/bin/python /opt/news_app/scripts/run_workers.py --queue content --worker-slot %(process_num)s --stats-interval 60"
directory=/opt/news_app
user=newsapp
autostart=true
autorestart=true
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
stopwaitsecs=900
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/news_app/workers.log
stderr_logfile=/var/log/news_app/workers.err.log
redirect_stderr=false
environment=ENVIRONMENT="production"

[program:news_app_workers_image]
command=/bin/bash -lc "/opt/news_app/.venv/bin/python /opt/news_app/scripts/run_workers.py --queue image --worker-slot 1 --stats-interval 60"
directory=/opt/news_app
user=newsapp
autostart=true
autorestart=true
stopwaitsecs=900
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/news_app/workers_image.log
stderr_logfile=/var/log/news_app/workers_image.err.log
redirect_stderr=false
environment=ENVIRONMENT="production"

[program:news_app_workers_media]
command=/bin/bash -lc "/opt/news_app/.venv/bin/python /opt/news_app/scripts/run_workers.py --queue media --worker-slot %(process_num)s --stats-interval 60"
directory=/opt/news_app
user=newsapp
autostart=true
autorestart=true
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
stopwaitsecs=900
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/news_app/workers_media.log
stderr_logfile=/var/log/news_app/workers_media.err.log
redirect_stderr=false
environment=ENVIRONMENT="production"

[program:news_app_workers_onboarding]
command=/bin/bash -lc "/opt/news_app/.venv/bin/python /opt/news_app/scripts/run_workers.py --queue onboarding --worker-slot 1 --stats-interval 60"
directory=/opt/news_app
user=newsapp
autostart=true
autorestart=true
stopwaitsecs=900
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/news_app/workers_onboarding.log
stderr_logfile=/var/log/news_app/workers_onboarding.err.log
redirect_stderr=false
environment=ENVIRONMENT="production"

[program:news_app_workers_twitter]
command=/bin/bash -lc "/opt/news_app/.venv/bin/python /opt/news_app/scripts/run_workers.py --queue twitter --worker-slot 1 --stats-interval 60"
directory=/opt/news_app
user=newsapp
autostart=true
autorestart=true
stopwaitsecs=900
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/news_app/workers_twitter.log
stderr_logfile=/var/log/news_app/workers_twitter.err.log
redirect_stderr=false
environment=ENVIRONMENT="production"

[program:news_app_workers_chat]
command=/bin/bash -lc "/opt/news_app/.venv/bin/python /opt/news_app/scripts/run_workers.py --queue chat --worker-slot 1 --stats-interval 60"
directory=/opt/news_app
user=newsapp
autostart=true
autorestart=true
stopwaitsecs=900
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/news_app/workers_chat.log
stderr_logfile=/var/log/news_app/workers_chat.err.log
redirect_stderr=false
environment=ENVIRONMENT="production"

[program:news_app_queue_watchdog]
command=/opt/news_app/scripts/start_queue_watchdog.sh
directory=/opt/news_app
user=newsapp
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/news_app/watchdog.log
stderr_logfile=/var/log/news_app/watchdog.err.log
redirect_stderr=false
environment=PATH="/opt/news_app/.venv/bin:/usr/local/bin:/usr/bin:/bin",QUEUE_WATCHDOG_MEDIA_STALE_HOURS="2",QUEUE_WATCHDOG_PROCESS_CONTENT_STALE_HOURS="2",QUEUE_WATCHDOG_ALERT_THRESHOLD="1"

[program:news_app_bgutil_provider]
command=/bin/bash -c "/opt/news_app/scripts/start_bgutil_provider.sh"
directory=/opt/news_app
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/news_app/bgutil_provider.log
stderr_logfile=/var/log/news_app/bgutil_provider.err.log
redirect_stderr=false
environment=BGUTIL_PROVIDER_IMAGE="brainicism/bgutil-ytdlp-pot-provider:1.3.1",BGUTIL_PROVIDER_HOST="127.0.0.1",BGUTIL_PROVIDER_PORT="4416",BGUTIL_PROVIDER_PULL_IMAGE="missing"

# Scrapers/discovery scheduling is managed by cron, not Supervisor.
# Keep only API + workers here.
#
# Recommended `newsapp` crontab entries:
# */15 * * * * flock -n /tmp/news_app_scrapers.lock /bin/bash -lc 'cd /opt/news_app && /opt/news_app/.venv/bin/python scripts/run_scrapers.py --show-stats --scrapers HackerNews Reddit Substack Techmeme Podcast Atom' >> /var/log/news_app/scrapers-cron.log 2>&1
# */15 * * * * flock -n /tmp/news_app_twitter.lock /bin/bash -lc 'cd /opt/news_app && /opt/news_app/.venv/bin/python scripts/run_twitter.py' >> /var/log/news_app/twitter.log 2>&1
# 0 3 * * 1 flock -n /tmp/news_app_feed_discovery.lock /bin/bash -lc 'cd /opt/news_app && /opt/news_app/.venv/bin/python scripts/run_feed_discovery.py' >> /var/log/news_app/feed-discovery-cron.log 2>&1
