# supervisor.conf
# News App Supervisor Configuration
# Copy to: /etc/supervisor/conf.d/news_app.conf
# Then run: supervisorctl reread && supervisorctl update

[program:news_app_server]
command=/opt/news_app/scripts/start_server.sh
directory=/opt/news_app
user=newsapp
autostart=true
autorestart=true
stopsignal=TERM
startretries=3
stdout_logfile=/var/log/news_app/server.log
stderr_logfile=/var/log/news_app/server.err.log
redirect_stderr=false
environment=PATH="/opt/news_app/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[program:news_app_workers_content]
command=/opt/news_app/scripts/run_workers.py --queue content --worker-slot %(process_num)s --stats-interval 60
directory=/opt/news_app
user=newsapp
autostart=true
autorestart=true
numprocs=2
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/news_app/workers.log
stderr_logfile=/var/log/news_app/workers.err.log
redirect_stderr=false
environment=PATH="/opt/news_app/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[program:news_app_workers_onboarding]
command=/opt/news_app/scripts/run_workers.py --queue onboarding --worker-slot 1 --stats-interval 60
directory=/opt/news_app
user=newsapp
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/news_app/workers_onboarding.log
stderr_logfile=/var/log/news_app/workers_onboarding.err.log
redirect_stderr=false
environment=PATH="/opt/news_app/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[program:news_app_workers_chat]
command=/opt/news_app/scripts/run_workers.py --queue chat --worker-slot 1 --stats-interval 60
directory=/opt/news_app
user=newsapp
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/news_app/workers_chat.log
stderr_logfile=/var/log/news_app/workers_chat.err.log
redirect_stderr=false
environment=PATH="/opt/news_app/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# Scrapers/discovery scheduling is managed by cron, not Supervisor.
# Keep only API + workers here.
#
# Recommended `newsapp` crontab entries:
# */15 * * * * cd /opt/news_app && /opt/news_app/.venv/bin/python scripts/run_scrapers.py --show-stats >> /var/log/news_app/scrapers-cron.log 2>&1
# 0 * * * * cd /opt/news_app && /opt/news_app/.venv/bin/python scripts/run_daily_news_digest.py >> /var/log/news_app/daily-news-digest.log 2>&1
# 0 3 * * 1 cd /opt/news_app && /opt/news_app/.venv/bin/python scripts/run_feed_discovery.py >> /var/log/news_app/feed-discovery-cron.log 2>&1
